# backbone/devtools/rootfs/etc/services.d/ghenghis/run
#!/usr/bin/with-contenv sh

exec /usr/bin/genghisapp --no-launch  --servers mongo:27017


