# backbone/base/rootfs/etc/services.d/syslog/run
#!/bin/execlineb -P

exec -c
fdmove -c 2 1

/sbin/syslogd -n -O -
