# application/mongodb/rootfs/etc/services.d/mongodb/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid mongodb mongod --dbpath /data/db
#exec /usr/bin/mongodb
