# backbone/url-monitor/rootfs/etc/services.d/url-monitor/run
#!/usr/bin/with-contenv sh

exec go run main.go -c config.yml;
