# docs/hugo/rootfs/etc/services.d/hugo/run
#!/usr/bin/with-contenv sh

exec /usr/local/sbin/hugo server   --source="/app" --destination="/app/public"
