# docs/hugo/rootfs/etc/services.d/hugo/run
#!/usr/bin/with-contenv sh

exec hugo server --source="/app" --baseURL=https://tracker.learn.cisco/docs --theme=hugo-theme-learn  --port=443 --bind=0.0.0.0

