# docs/hugo/rootfs/etc/services.d/hugo/run
#!/usr/bin/with-contenv sh

exec hugo version
