# docs/hugo/rootfs/etc/services.d/hugo/run
#!/usr/bin/with-contenv sh

exec hugo server --source="/app" --theme=hugo-theme-learn --port=80 --bind=0.0.0.0

