# docs/hugo/rootfs/etc/services.d/hugo/run
#!/usr/bin/with-contenv sh

hugo --source="/app" --destination="/app/public"
exec hugo server --source="/app" --destination="/app/public"
