# tools/ops/security/clamav/rootfs/etc/services.d/clamd/run
#!/usr/bin/env sh
exec clamd  -O /dev/stdout
