# tools/ops/security/clamav/rootfs/etc/services.d/freshclam/run
#!/usr/bin/env sh
exec freshclam -d  -O /dev/stdout
