# mariadb/rootfs/etc/services.d/mariadb/run
#!/bin/execlineb -P

foreground {
  s6-svwait -u /var/run/s6/services/syslog
}

exec -c
fdmove -c 2 1

/usr/bin/mysqld_safe  --socket=/var/lib/mysql/mysql.sock --user=daspanel --console --port=3306 --syslog
